# docker/root/etc/services.d/cloudflare-ddns/run
#!/usr/bin/with-contenv sh

/app/cloudflare.sh && crond -f
